# ansible/roles/dnsmasq/files/avahi.path
[Unit]
Description=My Avahi Service Path
Before=paths.target
RefuseManualStart=true
RefuseManualStop=true

[Path]
DirectoryNotEmpty=/run/systemd/avahi.dnsmasqd
MakeDirectory=yes
PathModified=/run/systemd/avahi.dnsmasqd
