# ansible/roles/dnsmasq/files/avahi.timer
[Unit]
Description=My Avahi Timer

[Timer]
OnBootSec=15s
OnUnitActiveSec=15s

[Install]
WantedBy=basic.target

